# config/kong/supervisor-pushpin.conf
[program:supervisor-pushpin]
command=/bin/bash sudo pushpin
autostart=true
autorestart=true
startretries=10
user=ideam
redirect_stderr=true
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
