# ubuntu/config/idps/log_parser.conf
[program:log_parser]
command=
autostart=true
autorestart=true
stderr_logfile=/var/log/log_parser.log
stdout_logfile=/var/log/log_parser.out.log